# configs/soc-mirror.service
[Unit]
Description=SOC Lab Traffic Mirror (boot-time catch-up reconcile)
Documentation=file:///var/home/Joost/Homelab/docs/guides/security_onion_setup.md
After=libvirtd.service
Requires=libvirtd.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/soc-mirror.sh --reconcile
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
